# ai-system-assistant/systemd/ai-assistant.service
[Unit]
Description=AI System Assistant - Local AI for Ubuntu System Management
After=network.target ollama.service
Wants=ollama.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/ai-system-assistant
ExecStart=/usr/bin/python3 -m src.cli
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ai-assistant

# Security hardening
NoNewPrivileges=false
ProtectSystem=full
ProtectHome=read-only
PrivateTmp=true

# Environment
Environment=PYTHONUNBUFFERED=1
Environment=AI_ASSISTANT_LOG=/var/log/ai-assistant.log

[Install]
WantedBy=multi-user.target
